# infra/systemd/trenches-executor.service
[Unit]
Description=Trenches Executor
After=network.target

[Service]
WorkingDirectory=%h/trenches
EnvironmentFile=%h/trenches/.env
ExecStart=/usr/bin/node services/executor/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

