# infra/systemd/trenches-policy-engine.service
[Unit]
Description=Trenches Policy Engine
After=network.target

[Service]
WorkingDirectory=%h/trenches
EnvironmentFile=%h/trenches/.env
ExecStart=/usr/bin/node services/policy-engine/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

