# infra/systemd/trenches-narrative-miner.service
[Unit]
Description=Trenches Narrative Miner
After=network.target

[Service]
WorkingDirectory=%h/trenches
EnvironmentFile=%h/trenches/.env
ExecStart=/usr/bin/node services/narrative-miner/dist/index.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

